# files/etc/crontabs/root
30 02 * * 1 sleep 5 && touch /etc/banner && reboot
0 2 * * * /usr/share/UnblockNeteaseMusic/update_core.sh 2>&1
0 2 * * * /usr/share/shadowsocksr/ssrplusupdate.sh >/var/log/ssrplus.log
